# server/supervisord/supervisord.conf
# -*- conf -*-

# [include]
# files = *.supervisor

[supervisord]

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[unix_http_server]
file = /var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:instantfm]
process_name = instantfm-%(process_num)s
directory=/home/feross/www/dev.instant.fm/publish/server/
command = python server.py --port=%(process_num)s
# Increase numprocs to run multiple processes on different ports.
numprocs = 4
numprocs_start = 7000